# project/scripts/package/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/ddg-idd-driver

update-rc.d ddg-idd-driver defaults

service ddg-idd-driver start