# package/DEBIAN/postrm
#!/bin/bash

update-rc.d ddg-idd-driver remove

#rm -fr /etc/ddg_idd_driver

